# platforms/macos/com.chip.zero-cost-warp-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chip.zero-cost-warp-guard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/chip/.config/zero-cost/warp-guard.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/chip/.config/zero-cost/warp-guard.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chip/.config/zero-cost/warp-guard.err.log</string>
</dict>
</plist>
